# services/agent_gateway/deploy/launchd/ai.kakeya.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Agent gateway LaunchAgent template (HEAD Mac only). Substitute __HOME__ and
     __WAN_WORKERS__ (comma list of all worker addrs). Add an AGENT_GATEWAY_API_KEY
     EnvironmentVariables entry to re-require X-API-Key. Install to ~/Library/LaunchAgents/. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.kakeya.gateway</string>
  <key>ProgramArguments</key>
  <array>
    <string>__HOME__/.venv-distwan/bin/python</string>
    <string>services/agent_gateway/server.py</string>
    <string>--host</string><string>127.0.0.1</string>
    <string>--port</string><string>8088</string>
  </array>
  <key>WorkingDirectory</key><string>__HOME__/openmontage-mac</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>__HOME__/.venv-distwan/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>WAN_WORKERS</key><string>__WAN_WORKERS__</string>
    <!-- Mode (with >1 worker):
         PIPELINE (default, no key): ONE job spans both Macs (head=proposer framework worker,
           headless=refine worker) via the orchestrator's --single-refine path. Higher res/quality.
         POOL: add  <key>AGENT_GATEWAY_WORKER_POOL</key><string>1</string>  -> each job runs DIRECT
           on ONE Mac, N jobs in parallel = N× throughput (every Mac worker must advertise framework). -->
    <key>ORCHESTRATOR_PATH</key><string>__HOME__/openmontage-mac/services/distributed_wan/grpc_orchestrator.py</string>
    <key>AGENT_GATEWAY_JOBS_DIR</key><string>__HOME__/.openmontage-jobs</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>__HOME__/.openmontage-logs/gateway.log</string>
  <key>StandardErrorPath</key><string>__HOME__/.openmontage-logs/gateway.log</string>
</dict>
</plist>
