# services/agent_gateway/deploy/launchd/ai.kakeya.mlxworker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- MLX worker LaunchAgent template. Substitute __HOME__ (e.g. /Users/allen),
     __HOST__ (127.0.0.1 on the head Mac, the LAN/Thunderbolt IP on a peer),
     __MODEL_DIR__ (e.g. __HOME__/wan21_mlx), and __MLX_OPS__ = the worker's ROLE:
       framework  -> proposer (head Mac, low-res T2V)
       refine     -> refiner  (headless Mac; SR upscale, or V2V if MLX_V2V_FLAG set)
       framework,refine -> both (single-Mac coarse-to-fine)
     Install to ~/Library/LaunchAgents/. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.kakeya.mlxworker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-dimsu</string>
    <string>__HOME__/.venv-distwan/bin/python</string>
    <string>grpc_worker.py</string>
    <string>--backend</string><string>mlx</string>
    <string>--host</string><string>__HOST__</string>
    <string>--port</string><string>50051</string>
    <string>--mlx-model-dir</string><string>__MODEL_DIR__</string>
    <string>--mlx-ops</string><string>__MLX_OPS__</string>
  </array>
  <key>WorkingDirectory</key><string>__HOME__/openmontage-mac/services/distributed_wan</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>__HOME__/.venv-distwan/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>MLX_TILING</key><string>aggressive</string>
    <key>MLX_RELATIVE_SPEED</key><string>0.12</string>
    <!-- Refiner only: MLX_SR_SHARPEN sets the SR unsharp percent (0 disables). If your
         mlx-video build has vid2vid, set MLX_V2V_FLAG (e.g. --video) for generative refine. -->
    <key>MLX_SR_SHARPEN</key><string>60</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ProcessType</key><string>Interactive</string>
  <key>StandardOutPath</key><string>__HOME__/.openmontage-logs/worker.log</string>
  <key>StandardErrorPath</key><string>__HOME__/.openmontage-logs/worker.log</string>
</dict>
</plist>
